# files/cpu-frequency/omen-mode.service
[Unit]
Description=Set Hardware Performance

[Service]
ExecStart=/home/skeletony007/.local/scripts/omen-mode

[Install]
WantedBy=multi-user.target
